# automine.service
[Unit]
Description=XCB Autominer2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/miner
ExecStart=/bin/bash /path/to/miner/automine.sh
Restart=always
RestartSec=10
KillMode=control-group
TimeoutStopSec=10
Environment=TERM=xterm
StandardOutput=append:/var/log/autominer2/autominer2.log
StandardError=append:/var/log/autominer2/autominer2.log

# Optionaler eigener Hostname fürs Dashboard
# Environment=CUSTOM_HOSTNAME=miner-berlin-01

[Install]
WantedBy=multi-user.target
